# ui-verifiers.service
[Unit]
Description=UI Verifiers
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/ui-verifiers
ExecStart=/bin/bash -lc "uv run uvicorn ui_verifiers.server:app --host 0.0.0.0"
Restart=always

[Install]
WantedBy=default.target